# schedule/com.linkedin.autoposter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.linkedin.autoposter</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/miniconda3/envs/linkedin-autoposter/bin/python</string>
    <string>-m</string>
    <string>src.main</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/dhamotharanpalanisamy/projects/linkedin-autoposter</string>

  <!-- Run daily at 08:30 local time -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <!-- Catch up if machine was asleep at scheduled time -->
  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/dhamotharanpalanisamy/Library/Logs/linkedin-autoposter.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/dhamotharanpalanisamy/Library/Logs/linkedin-autoposter-err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/miniconda3/envs/linkedin-autoposter/bin:/opt/miniconda3/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/dhamotharanpalanisamy</string>
  </dict>
</dict>
</plist>
